# systemd/stone.target
[Unit]
Description=stone target allowing to start/stop all stone*@.service instances at once

[Install]
WantedBy=multi-user.target
