# systemd/stone-mds@.service.in
[Unit]
Description=Stone metadata server daemon
PartOf=stone-mds.target
After=network-online.target local-fs.target time-sync.target
Before=remote-fs-pre.target stone-mds.target
Wants=network-online.target local-fs.target time-sync.target remote-fs-pre.target stone-mds.target

[Service]
Environment=CLUSTER=stone
EnvironmentFile=-@SYSTEMD_ENV_FILE@
ExecReload=/bin/kill -HUP $MAINPID
ExecStart=/usr/bin/stone-mds -f --cluster ${CLUSTER} --id %i --setuser stone --setgroup stone
LimitNOFILE=1048576
LimitNPROC=1048576
LockPersonality=true
MemoryDenyWriteExecute=true
NoNewPrivileges=true
PrivateDevices=yes
PrivateTmp=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=full
Restart=on-failure
RestrictSUIDSGID=true
StartLimitBurst=3
StartLimitInterval=30min
TasksMax=infinity

[Install]
WantedBy=stone-mds.target
