# debian/stone-base.prerm
#!/bin/sh
# vim: set noet ts=8:

set -e

case "$1" in
    remove)
	[ -x /sbin/stop ] && stop stone-all || true
	invoke-rc.d stone stop || {
	    RESULT=$?
	    if [ $RESULT != 100 ]; then
		exit $RESULT
	    fi
	}
	;;

    *)
	;;
esac

#DEBHELPER#

exit 0
